# scripts/habit-tracker-maintenance.service
[Unit]
Description=Tracker Database Maintenance
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/habit-tracker
ExecStart=/opt/habit-tracker/scripts/db-maintenance.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
